# docker/supervisor/asciinema.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true

[program:rails]
command=bundle exec rails server
redirect_stderr=true
directory=/app